# services/yukarin-voice-caster-sleep-preparation.timer
[Unit]
Description=Play yukarin sleep preparation voice

[Timer]
OnCalendar=*-*-* 13:00

[Install]
WantedBy=timers.target
