# services/yukarin-voice-caster-sleep-preparation.service
[Unit]
Description=Play yukarin sleep preparation voice

[Service]
Type=oneshot
ExecStart=%h/yukarin-voice-caster/main.py sleep-preparation.mp3
