# services/yukarin-voice-caster-melatonin.timer
[Unit]
Description=Play yukarin melatonin voice

[Timer]
OnCalendar=*-*-* 12:00

[Install]
WantedBy=timers.target
